# supervisord.conf
[supervisord]
nodaemon=true

[program:AppSettings]
command=python3 /opt/AppSettings/app.py
autostart=true
autorestart=true